# root/etc/s6-overlay/s6-rc.d/init-bazel/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

if [[ -d /config/clwb ]]; then
	rm -rf /config/clwb/bazel-*
fi

if [[ -d /config/workspace ]]; then
	rm -rf /config/workspace/bazel-*
fi
