# root/config/.config/autostart/idea_backend.desktop
[Desktop Entry]
Exec=konsole -e /defaults/idea_backend.sh
Icon=
Name=IntelliJ Backend
Path=
Terminal=False
Type=Application
