# root/etc/s6-overlay/s6-rc.d/svc-agentlib-forwarder/run
#! /usr/bin/execlineb -P

# Redirect stderr to stdout.
fdmove -c 2 1

# Forward port 0.0.0.0:5006 to 127.0.0.1:5005
socat TCP4-LISTEN:5006,fork,reuseaddr TCP4:127.0.0.1:5005
